# build/deb/stable/DEBIAN/postinst
#!/bin/sh

adduser --system --disabled-password opensemanticetl

# start on boot
update-rc.d opensemanticsearch defaults
update-rc.d opensemanticsearch-filemonitoring defaults

# start workers after installation (but not the filemonitoring, since no pahts configured yet)
service opensemanticsearch start
