# build/deb/stable/DEBIAN/prerm
#!/bin/sh

systemctl disable opensemanticetl-filemonitoring

systemctl stop opensemanticetl-filemonitoring

systemctl disable opensemanticetl

systemctl stop opensemanticetl

exit 0
